# infra/cron/booking.crontab
# Booking scheduler sidecar jobs (Europe/Lisbon via container TZ).
# Complete approved bookings whose end time has passed.
0 1 * * * python manage.py complete_bookings
